# deploy/systemd/nsqlokupd@.service
[Unit]
Description= nsqlookupd service file generated by Enso
Documentation=man:nsqlookupd(1)
After=network.target

[Service]
User=cc
Group=cc
ExecStart=/opt/nsq-latest/bin/nsqlookupd 
Restart=on-failure
WorkingDirectory=/opt/nsq-latest/bin
StandardOutput=null
StandardError=null
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
